# nodejs-wsd
#!/lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          nodejs-wsd
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: NodeJS WebSocket Server
### END INIT INFO

# Author: Benjamin Kellermann <benjamin.kellermann@mgm-sp.com>
DESC="NodeJS WebSocket Server"
DAEMON=/var/www/nodejs-ws
START_ARGS="--background --make-pidfile"
STOP_ARGS="--remove-pidfile"
NAME="node"
do_stop_override() {
				start-stop-daemon --stop --retry=TERM/30/KILL/5 \
						$STOP_ARGS \
						${PIDFILE:+--pidfile ${PIDFILE}} --name $NAME
				return $RETVAL
}
